# config-examples/eve_tail2duckdb.service
[Unit]
Description=Suricata eve.json -> DuckDB ingester
After=suricata.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER
Environment=PATH=/home/YOUR_USER/.local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /home/YOUR_USER/eve_tail2duckdb.py \
    --eve /var/log/suricata/eve.json \
    --db /home/YOUR_USER/eve.duckdb \
    --batch 50000 --read-gap 1.0 --reopen-every 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
